# containerfiles/etc/supervisor/conf.d/xrdp.conf
[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
user=root
priority=10
autorestart=true

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
user=root
priority=11
autorestart=true
